# .config/autostart/udiskie.desktop
[Desktop Entry]
Exec=/usr/bin/udiskie --tray
Name=udiskie
Comment=automounter for removable media
Type=Application
